# crontab
*/3 * * * * sh -c "cd ~/git/monero.fail; ./bin/cmd validate"
30 */4 * * * sh -c "cd ~/git/monero.fail; ./bin/cmd check"
0 0 * * * sh -c "cd ~/monero.fail; ./manage.sh export"
*/5 * * * * sh -c "cd ~/git/monero.fail; ./bin/cmd get_peers"
